# wol-proxy.service
# /etc/systemd/system/wol-proxy.service
# sudo cp wol-proxy.service  /etc/systemd/system/wol-proxy.service
# sudo apt-get install -y python3-redis
# sudo systemctl daemon-reload
# sudo systemctl enable --now wol-proxy
# sudo journalctl -u wol-proxy -f

[Unit]
Description=WOL proxy (Redis pubsub -> UDP broadcast)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/env python3 /usr/local/bin/wol-proxy.py
Restart=always
RestartSec=1
Environment=REDIS_HOST=swarm1
Environment=REDIS_PORT=6379
Environment=REDIS_DB=1
Environment=WOL_CHANNEL=wol:wake
# Run as root so it can broadcast without restrictions
User=root

[Install]
WantedBy=multi-user.target
